# systemd/memstream-distill.timer
[Unit]
Description=memstream distill — hourly
Requires=memstream-distill.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
